# script/bank-of-asgard-agent.service
[Unit]
Description=Bank of Asgard Transactions Agent (WebSocket)
Documentation=https://github.com/asgardeo/bank-of-asgard
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/bank-of-asgard/transactions-agent

# Activates the venv and starts uvicorn on all interfaces — DO LB accesses via VM's IP
ExecStart=/bin/bash -lc 'source %h/bank-of-asgard/transactions-agent/.venv/bin/activate && uvicorn app.service:app --host 0.0.0.0 --port 8011'

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bank-of-asgard-agent

[Install]
WantedBy=default.target
