# script/bank-of-asgard-agent.service
[Unit]
Description=Bank of Asgard Transactions Agent (WebSocket)
Documentation=https://github.com/asgardeo/bank-of-asgard
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=boa
Group=boa
WorkingDirectory=/home/boa/bank-of-asgard/transactions-agent

ExecStart=/bin/bash -lc 'uvicorn app.service:app --host 0.0.0.0 --port 8011'

Restart=on-failure
RestartSec=5


StandardOutput=journal
StandardError=journal
SyslogIdentifier=bank-of-asgard-agent

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/boa/bank-of-asgard

[Install]
WantedBy=default.target
